# deploy-as-code/dockerfiles/redoc/crontab
0 0-23/4 * * * root /usr/bin/python /opt/eGov/update_redoc_html.py > /tmp/cron.log 2>&1


